# scripts/startup/medical-ai-service.service
[Unit]
Description=Medical AI Service for DocuSeal Integration
After=network.target
Wants=network.target

[Service]
Type=simple
User=rvalen
Group=rvalen
WorkingDirectory=/home/rvalen/Projects/Milestones/wound-care-react-laravel/scripts
ExecStart=/home/rvalen/Projects/Milestones/wound-care-react-laravel/scripts/venv/bin/python medical_ai_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medical-ai-service

# Environment variables
Environment=PYTHONPATH=/home/rvalen/Projects/Milestones/wound-care-react-laravel/scripts
Environment=API_HOST=0.0.0.0
Environment=API_PORT=8080
Environment=ENABLE_LOCAL_FALLBACK=true

# Load additional environment variables from Laravel .env file
EnvironmentFile=-/home/rvalen/Projects/Milestones/wound-care-react-laravel/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/rvalen/Projects/Milestones/wound-care-react-laravel/storage/logs

[Install]
WantedBy=multi-user.target 